# config_templates/systemd/sing-box.service.template
# ======================================================================
# Systemd Service Unit File Template for Sing-Box
# Generated by: deploy_proxy_platform.sh
# Purpose: Run Sing-Box proxy backend as a non-root service.
# Location: /etc/systemd/system/sing-box.service
# ======================================================================

[Unit]
Description=Sing-Box Service
Documentation=https://sing-box.sagernet.org
# Ensure network is up and DNS resolution is available before starting
After=network.target nss-lookup.target

[Service]
# --- Execution ---
# User and Group to run the service as (non-root)
User=${SINGBOX_USER}
Group=${SINGBOX_GROUP}
# Working directory where config.json resides
WorkingDirectory=${SINGBOX_CONFIG_DIR}
# Command to start Sing-Box, pointing to the correct config file
ExecStart=${SINGBOX_INSTALL_DIR}/sing-box run -c ${SINGBOX_CONFIG_DIR}/config.json

# --- Process Management ---
# Restart the service if it fails
Restart=on-failure
# Time to wait before restarting
RestartSec=10s
# Increase the number of allowed open file descriptors
LimitNOFILE=infinity

# --- Security Hardening ---
# Prevent the service process from gaining new privileges
NoNewPrivileges=true
# Make system directories read-only for the service
ProtectSystem=strict
# Make user home directories inaccessible
ProtectHome=true
# Use a private temporary directory
PrivateTmp=true
# Restrict access to device nodes
PrivateDevices=true
# Clear default capabilities unless specific ones are proven necessary.
# For listening on loopback high ports, usually none are needed.
CapabilityBoundingSet=
AmbientCapabilities=
# Optionally restrict network access further if needed, e.g., only allow loopback
# IPAddressAllow=127.0.0.1 ::1
# IPAddressDeny=any

[Install]
# Enable the service to start on boot in the multi-user target
WantedBy=multi-user.target
