# old_shit/config_templates/systemd/subscription-app.service.template
# ======================================================================
# Systemd Service Unit File Template for Subscription Flask App (Gunicorn)
# Generated by: deploy_proxy_platform.sh
# Purpose: Run the Python Flask web application as a non-root service.
# Location: /etc/systemd/system/subscription-app.service
# ======================================================================

[Unit]
Description=Subscription Flask App Service (Gunicorn)
# Ensure network is available before starting
After=network.target

[Service]
# --- Execution ---
# User and Group to run the service as (non-root)
User=${SUBAPP_USER}
Group=${SUBAPP_GROUP}
# Set the working directory to the application's root
WorkingDirectory=${SUBSCRIPTION_APP_DIR}

# --- Environment Variables ---
# Pass necessary configuration from the deployment script to the Flask app
# Note: Sensitive information like passwords/UUIDs are included here.
# Consider alternative secret management for higher security needs.
Environment="PATH=${SUBSCRIPTION_APP_DIR}/venv/bin" # Ensure venv PATH is available
Environment="MAIN_DOMAIN=${MAIN_DOMAIN}"
Environment="SUBSCRIPTION_DOMAIN=${SUBSCRIPTION_DOMAIN}"
Environment="SUBSCRIPTION_BASE64_PATH=${SUBSCRIPTION_BASE64_PATH}"
Environment="API_BASE64_PATH_PREFIX=${API_BASE64_PATH_PREFIX}"
Environment="VLESS_PORT=${VLESS_HTTPUPGRADE_PORT}"
Environment="VLESS_UUID=${VLESS_UUID}" # Initial VLESS UUID
Environment="VLESS_PATH=${VLESS_PATH}"
Environment="HYSTERIA2_PORT=${HYSTERIA2_PORT}"
Environment="HYSTERIA2_PASSWORD=${HYSTERIA2_PASSWORD}"
Environment="FLASK_APP=app.py" # Name of the main Flask application file

# --- Command ---
# Execute Gunicorn WSGI server from the virtual environment
# - Bind to loopback interface only (HAProxy provides external access)
# - Specify number of worker processes (adjust based on server resources, e.g., 2*cores+1)
# - Log access and errors to stdout/stderr so journald can capture them
ExecStart=${SUBSCRIPTION_APP_DIR}/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:${SUBSCRIPTION_APP_LISTEN_PORT} \
    --access-logfile - \
    --error-logfile - \
    app:app

# --- Process Management ---
Restart=always          # Restart the service if it crashes
RestartSec=5s           # Wait 5 seconds before restarting
StandardOutput=journal  # Redirect stdout to systemd-journald
StandardError=journal   # Redirect stderr to systemd-journald

# --- Security Hardening ---
NoNewPrivileges=true            # Prevent gaining additional privileges
ProtectSystem=full              # Make /usr, /boot, /etc largely read-only
ProtectHome=true                # Make user home directories inaccessible
PrivateTmp=true                 # Use a private /tmp and /var/tmp
PrivateDevices=true             # Restrict device access in /dev
ReadWritePaths=${SUBSCRIPTION_APP_DIR}/ # Explicitly allow writing within the app dir (e.g., for cache/logs if needed)
CapabilityBoundingSet=          # Clear default capabilities
AmbientCapabilities=            # Clear default capabilities

[Install]
# Enable the service to start automatically on boot
WantedBy=multi-user.target
